# dist/moraga-dj.service
[Unit]
Description=Open Source Event Manager Delayed Jobs
Wants=moraga.service
Before=moraga.service

[Service]
User=moraga
EnvironmentFile=/home/moraga/events/current/.env.production
WorkingDirectory=/home/moraga/events/current
ExecStart = /usr/bin/bundle exec bin/delayed_job -e production start
ExecStop = /usr/bin/bundle exec bin/delayed_job -e production stop
Type = forking
PIDFile = /home/moraga/events/current/tmp/pids/delayed_job.pid

[Install]
WantedBy=multi-user.target

